# wifi_setup/nuggies-wifi-setup.service
[Unit]
Description=Nuggies WiFi Setup / Captive Portal
Documentation=https://github.com/nuggies/nuggies_pi_displays
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=simple
ExecStart=/bin/bash __PROJECT_DIR__/scripts/wifi_setup_service.sh
WorkingDirectory=__PROJECT_DIR__
User=root
StandardOutput=journal
StandardError=journal
Restart=no

[Install]
WantedBy=multi-user.target
